# packages/collective.funkbot/collective.funkbot-1.0.zip/collective.funkbot-1.0/collective/funkbot/templates/funkbot/etc/supervisord.conf
[unix_http_server]
file=%(here)s/../var/supervisor.sock
chmod=0600

[supervisorctl]
serverurl=unix://%(here)s/../var/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
logfile=%(here)s/../var/log/supervisord.log
logfile_maxbytes=5MB
logfile_backups=10
loglevel=info
pidfile=%(here)s/../var/supervisord.pid ;
childlogdir=%(here)s/../var/log
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)
directory=%(here)s

[program:zeo]
command = %(here)s/../bin/zeo fg
autostart= true
autorestart = true
startsecs = 10
priority = 100

[program:zeo-instance1]
command = %(here)s/../bin/zeo-instance1 console
startsecs = 60
priority = 2
redirect_stderr = true
autostart= true
autorestart = true
priority = 500

[program:zeo-instance2]
command = %(here)s/../bin/zeo-instance2 console
startsecs = 60
priority = 2
redirect_stderr = true
autostart= true
autorestart = true
priority = 500

[groups]
programs = zeo-instance1 zeo-instance2

[program:pound]
command = %(here)s/../bin/pidproxy %(here)s/../parts/pound/var/pound.pid %(here)s/../bin/poundrun
autostart= true
autorestart = true
priority = 1000

[program:monitor]
command = %(here)s/../bin/fl-monitor-ctl %(here)s/../etc/monitor.conf startd
autostart= true
autorestart = unexpected
exitcodes = 1
priority = 2000
