# packages/aditam.agent/aditam.agent-0.1.tar.gz/aditam.agent-0.1/debian/python-aditamagent.aditam-agent.init
#!/bin/sh 
#
### BEGIN INIT INFO
# Provides:          aditam.agent
# Required-Start:    $network $local_fs
# Required-Stop:     
# Should-Start:      $named
# Should-Stop:       
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Automated and DIstributed TAsk Manager agent
### END INIT INFO

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

NAME="aditam-agent"              # Introduce the short server's name here
PIDFILE=/var/run/$NAME.pid 
DAEMON=/usr/bin/aditam-agent.py     # Introduce the server's location here
DESC=""                           # Introduce a short description here


# Include defaults if available
if [ -f /etc/default/$NAME ] ; then
	. /etc/default/$NAME
fi

if [ "$ENABLED" = "0" ]; then
    echo "$NAME disabled, see /etc/default/aditam-agent"
    exit 0
fi

DAEMON_OPTS="--pidfile=$PIDFILE $OPTIONS"

start() {
	$DAEMON $DAEMON_OPTS start
	return $?
}

stop() {
	$DAEMON $DAEMON_OPTS stop
	return $?
}

case "$1" in
  start)
	echo "Starting $DESC " "$NAME"
	start
	;;
  stop)
	stop
        rm -f $PIDFILE
        ;;
  restart|force-reload)
        echo -n "Restarting $DESC" "$NAME"
        stop
        # Wait some sensible amount, some server need this
        [ -n "$DIETIME" ] && sleep $DIETIME
        start
        [ -n "$STARTTIME" ] && sleep $STARTTIME
	;;
  *)
	N=/etc/init.d/$NAME
	echo "Usage: $N {start|stop|restart|force-reload}" >&2
	exit 1
	;;
esac

exit 0
