# packages/auf-refer/auf-refer-1.1.1.tar.gz/auf-refer-1.1.1/debian/cron.d
0 7,13	* * *	auf-refer	[ -x /usr/sbin/auf-refer ] && /usr/sbin/auf-refer -u
