# packages/wheeljack/wheeljack-0.4.tar.gz/wheeljack-0.4/examples/buildout-install/etc/supervisord.conf
[unix_http_server]
file=%(here)s/../var/supervisor.sock
chmod=0600

[supervisorctl]
serverurl=unix://%(here)s/../var/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
logfile=%(here)s/../var/log/supervisord.log
logfile_maxbytes=5MB
logfile_backups=10
loglevel=info
pidfile=%(here)s/../var/supervisord.pid ;
childlogdir=%(here)s/../var/log
directory=%(here)s

[program:server]
command=%(here)s/../bin/django pasteserver --port=8000 --host=0.0.0.0
priority=10

[program:builder]
command=%(here)s/../bin/wheeljack continuous --settings=autobot.settings
environment=PYTHONPATH=%(here)s/..